# ansible/templates/ontreenode.service.j2
[Unit]
Description=OnTree Node Service
After=network.target

[Service]
User={{ app_user }}
Group={{ app_user }}
WorkingDirectory={{ app_path }}
Environment="PORT={{ app_port }}"
Environment="DATABASE_PATH={{ app_path }}/ontree.db"
ExecStart={{ app_path }}/{{ binary_name }}
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target