# ansible/ontreenode.service.j2
[Unit]
Description={{ app_name }} Go service
After=network.target

[Service]
Type=simple
User={{ app_user }}
Group={{ app_user }}
WorkingDirectory={{ app_path }}
ExecStart={{ app_path }}/{{ binary_name }}
Restart=always
RestartSec=10

# Environment variables
Environment="LISTEN_ADDR=:{{ app_port }}"
Environment="DB_PATH={{ app_path }}/ontree.db"
Environment="APPS_DIR={{ app_root_dir }}/apps"

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths={{ app_path }} {{ app_root_dir }}/apps

[Install]
WantedBy=multi-user.target